# debian/prerm
#!/bin/sh
set -e

# Stop service before removal
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl stop teleddns 2>/dev/null || true
    systemctl disable teleddns 2>/dev/null || true
fi
